# backend/deploy/husk-processor.cron
*/10 6-21 * * * /home/tkjelsrud/husk/backend/deploy/run_processor.sh
0 22 * * * /home/tkjelsrud/husk/backend/deploy/run_processor.sh
