# backend/deploy/husk-processor.cron
*/10 6-22 * * * cd /home/tkjelsrud/husk && ./.venv/bin/python -m backend.processor.main --once >> /home/tkjelsrud/husk/backend/.state/processor.log 2>&1
